# supervisor/conf.d/supervisor-starter.conf
[program:supervisor-starter]
process_name=%(program_name)s
command=python /etc/supervisor/supervisor_starter.py
autostart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
